# docker/bitcoin/regtest/crontab
* * * * * /scripts/mine.sh
* * * * * sleep 30 && /scripts/mine.sh
